# config/runit/haproxy/run
#!/bin/sh
exec haproxy -f /container/config/haproxy/haproxy.cfg 