# config/runit/openbox/run
#!/bin/sh
/bin/su -l -s /bin/sh -c "sleep 5; DISPLAY=:0 app" app &
exec /bin/su -l -s /bin/sh -c "export DISPLAY=:0; exec openbox" app